# tools/dind/config/services/openvpn@client.service.d/uuid.conf
[Service]
WorkingDirectory=/usr/src/app
EnvironmentFile=/usr/src/app/config/env
ExecStartPre=/usr/src/app/vpn-init
TimeoutStartSec=0
Restart=always
