# tools/dind/fake.service
[Unit]
Description=Fake service

[Service]
ExecStart=
ExecStart=/bin/true
