# openrc/node-supervisor
#!/sbin/openrc-run

description="Resin.io nodejs supervisor"
pid_file=/var/run/resin-supervisor-node.pid
waiting_time=10

start() {
	ebegin "Starting $description"
	start-stop-daemon --start --exec node src/app.js \
	--chdir /usr/src/app \
	--make-pidfile --pidfile "$pid_file" --background \
	--stdout /dev/console \
	--stderr /dev/console \
	--wait "$waiting_time"
	eend $?
}

stop() {
	ebegin "Stopping $description"
	start-stop-daemon --stop --exec node src/app.js --pidfile $pid_file --wait "$waiting_time"
	eend $?
}
