# tools/dind/config/services/docker.service.d/aufs.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon --storage-driver aufs --host=fd://
