# tools/dind/config/services/docker.service.d/aufs.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dind /usr/bin/docker daemon --storage-driver overlay2 --host=fd://
