# tools/dind/config/services/docker.service.d/overlayfs.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker -s overlay -d -H fd://
