# tools/dind/config/services/docker.service.d/overlayfs.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon --storage-driver overlay --host=fd://
